# automation/com.codex.epi-dossier.overnight-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.codex.epi-dossier.overnight-build</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/.venv/bin/python</string>
      <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/src/public_publish.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>EPI_DOSSIER_PUBLIC_PUBLISH_TIMEOUT_SECONDS</key>
      <string>2700</string>
      <key>EPI_DOSSIER_PUBLIC_PUBLISH_STALE_LOCK_SECONDS</key>
      <string>3600</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <array>
      <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>15</integer></dict>
      <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>15</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/logs/overnight-build.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/logs/overnight-build.err.log</string>
  </dict>
</plist>
