# automation/com.codex.epi-dossier.public-publish-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.codex.epi-dossier.public-publish-watchdog</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/.venv/bin/python</string>
      <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/src/public_publish_watchdog.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>EPI_DOSSIER_WATCHDOG_STALE_MINUTES</key>
      <string>45</string>
      <key>EPI_DOSSIER_WATCHDOG_TIMEOUT_SECONDS</key>
      <string>15</string>
      <key>EPI_DOSSIER_WATCHDOG_MANIFEST_URL</key>
      <string>https://dteichrow.github.io/app_exports/manifest.json</string>
      <key>EPI_DOSSIER_WATCHDOG_PUBLISH_MODE</key>
      <string>github_actions</string>
      <key>EPI_DOSSIER_WATCHDOG_GITHUB_REPOSITORY</key>
      <string>dteichrow/epi-dossier</string>
      <key>EPI_DOSSIER_WATCHDOG_GITHUB_WORKFLOW</key>
      <string>Newsdesk public publish</string>
      <key>EPI_DOSSIER_WATCHDOG_GH_BIN</key>
      <string>/opt/homebrew/bin/gh</string>
      <key>EPI_DOSSIER_WATCHDOG_REMOTE_DISPATCH_COOLDOWN_MINUTES</key>
      <string>45</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <array>
      <dict><key>Minute</key><integer>5</integer></dict>
      <dict><key>Minute</key><integer>15</integer></dict>
      <dict><key>Minute</key><integer>25</integer></dict>
      <dict><key>Minute</key><integer>35</integer></dict>
      <dict><key>Minute</key><integer>45</integer></dict>
      <dict><key>Minute</key><integer>55</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/logs/public-publish-watchdog.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/logs/public-publish-watchdog.err.log</string>
  </dict>
</plist>
