# automation/com.codex.epi-dossier.overnight-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.codex.epi-dossier.overnight-build</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/.venv/bin/python</string>
      <string>-m</string>
      <string>src.site_build</string>
      <string>--days</string>
      <string>7</string>
      <string>--output-mode</string>
      <string>both</string>
      <string>--deploy-dir</string>
      <string>docs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <array>
      <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/logs/overnight-build.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/devinteichrow/Downloads/Work and Statistics/Blogs/epi-dossier/logs/overnight-build.err.log</string>
  </dict>
</plist>
